# autostart/.config/autostart/Dunst.desktop
[Desktop Entry]
Categories=Utility;
Comment=Notification daemon
Name=Dunst
StartupNotify=true
Terminal=false
Type=Application
Exec=/usr/bin/dunst
