# autostart/.config/autostart/Picom.desktop
[Desktop Entry]
Categories=Utility;
Comment=Compositor
Name=Picom
StartupNotify=true
Terminal=false
Type=Application
Exec=/usr/bin/picom