# autostart/.config/autostart/Udiskie.desktop
[Desktop Entry]
Categories=Utility;
Comment=USB systray applet
Name=Udiskie
StartupNotify=true
Terminal=false
Type=Application
Exec=/usr/bin/udiskie
